# meta-mxxf1/recipes-support/history-daemon/files/history-daemon.service
[Unit]
Description=User history

[Service]
Type=forking
ExecStart=/usr/bin/historydaemon
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
